# addon-hyperhdr/rootfs/etc/services.d/hyperhdr/run
#!/usr/bin/with-contenv bash
set -e

mkdir -p /config/hyperhdr

exec /usr/bin/hyperhdr \
  --userdata /config/hyperhdr \
  --webPort 8090 \
  --logLevel INFO
