# addon-hyperhdr/rootfs/etc/services.d/hyperhdr/run
#!/usr/bin/with-contenv bash
exec /usr/bin/hyperhdr \
  --userdata /config/hyperhdr \
  --webPort 8090 \
  --logLevel info
